# 8.1/config/etc/crontab
*/15  *  *  *  * su www-data -c "php -f /var/www/owncloud/cron.php"
