# 7.0/config/etc/supervisor/conf.d/cron.conf
[program:cron]
command=/usr/sbin/cron -f
redirect_stderr=true
stopsignal=QUIT
autorestart=true
