# 7.0/config/etc/supervisor/conf.d/phpfpm.conf
[program:phpfpm]
command=/usr/sbin/php5-fpm
stopsignal=QUIT
redirect_stderr=true
autorestart=true
